# client/logcleaner/kcplogcleaner.service
[Unit]
Description=KCP log cleaner

[Service]
Type=simple
ExecStart=/etc/Cleaner/cleanKCPLog.sh



